# xiaobaozi.supervisor.conf
[program:xiaobaozi_robot]
command=/usr/bin/python /root/xiaobaozi_robot/run.py
stopasgroup=true
killasgroup=true
restart=always
user=root
